# root/etc/services.d/aria2/run
#!/usr/bin/with-contenv bash

echo 'Exiting aria2'